# cron.d/s3fs-cache-cleanup
0 4 */1 * * USERNAME /usr/bin/s3fs-delcache.sh S3BUCKET /var/tmp 1073741824
