# cron.d/s3fs-hk-watchdog
0 */1 * * * USERNAME /usr/bin/s3fs-hk-watchdog.sh
